# cloud/cloudx_cloud/data/systemd/cloudx-legacy-health-bridge.service
[Unit]
Description=Publish the migration-only Cloudx legacy health bridge
After=cloudx-health.service
ConditionPathExists=/run/cloudx/health.json

[Service]
Type=oneshot
EnvironmentFile=/etc/cloudx/legacy-health-bridge.env
ExecStart=/usr/bin/python3 ${CLOUDX_LEGACY_HEALTH_BRIDGE_ARTIFACT} legacy-health-bridge --source /run/cloudx/health.json --publish-to /var/lib/cloudx/health/v1.json
User=root
Group=root
UMask=0022
StateDirectory=cloudx/health
StateDirectoryMode=0755
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadOnlyPaths=/opt/cloudx/releases /etc/cloudx /run/cloudx
ReadWritePaths=/var/lib/cloudx/health
InaccessiblePaths=/var/lib/codex-gateway/cliproxy-auth /var/lib/codex-gateway/cliproxy-auth-archive /etc/cliproxy /etc/codex-gateway/import-client-keys
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX
