# cloud/cloudx_cloud/data/systemd/cloudx-cpa-health.service
[Unit]
Description=Probe and reversibly isolate Cloudx CPA accounts with the native signed runtime
After=network-online.target mihomo.service cliproxy.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
Environment=HTTP_PROXY=http://127.0.0.1:7890
Environment=HTTPS_PROXY=http://127.0.0.1:7890
Environment=CLOUDX_CPA_PROXY_URL=http://127.0.0.1:7890
Environment=CLOUDX_WARNING_AVAILABLE_ACCOUNTS=3
Environment=CLOUDX_FAILURE_CONFIRMATIONS=3
ExecStart=/usr/bin/python3 /opt/cloudx/current/cloudx-cloud.pyz cpa-health
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
StateDirectory=cloudx/cpa-health
StateDirectoryMode=0700
ReadOnlyPaths=/opt/cloudx/releases
ReadWritePaths=/var/lib/cloudx/cpa-health /var/lib/codex-gateway/cliproxy-auth /var/lib/codex-gateway/cliproxy-auth-archive /var/lib/codex-gateway/cliproxy-auth-failures
InaccessiblePaths=/etc/phi /etc/msmtprc /etc/codex-quota-monitor.env /etc/codex-quota-monitor.gmail-app-password
PrivateDevices=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
