# cloud/cloudx_cloud/data/systemd/cloudx-account-state.service
[Unit]
Description=Publish aggregate account state for active Cloudx health
After=cloudx-cpa-health.service

[Service]
Type=oneshot
EnvironmentFile=/etc/cloudx/cloudx-shadow.env
ExecStart=/usr/bin/env CLOUDX_ACCOUNT_STATE_PATH=/run/cloudx-account-state/accounts.json /usr/bin/python3 /opt/cloudx/current/cloudx-cloud.pyz adapt-account-state --json
User=root
Group=root
UMask=0022
RuntimeDirectory=cloudx-account-state
RuntimeDirectoryMode=0755
RuntimeDirectoryPreserve=yes
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadOnlyPaths=/opt/cloudx/releases /etc/cloudx /var/lib/cloudx/cpa-health/state.json
ReadWritePaths=/run/cloudx-account-state
InaccessiblePaths=/var/lib/codex-gateway/cliproxy-auth /var/lib/codex-gateway/cliproxy-auth-archive /etc/cliproxy /etc/codex-gateway/import-client-keys
LockPersonality=true
MemoryDenyWriteExecute=true
