# cloud/systemd/cloudx-shadow-health.timer
[Unit]
Description=Publish Cloudx shadow health every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=20s
Persistent=true

[Install]
WantedBy=timers.target
