# cloud/systemd/cloudx-shadow-account-state.service
[Unit]
Description=Adapt aggregate credential health for Cloudx shadow health
After=cloudx-cpa-health.service

[Service]
Type=oneshot
EnvironmentFile=/etc/cloudx/cloudx-shadow.env
ExecStart=/usr/bin/python3 ${CLOUDX_CLOUD_ARTIFACT} adapt-account-state --json
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/opt/cloudx/releases /etc/cloudx /var/lib/cloudx/cpa-health/state.json
ReadWritePaths=/run/cloudx-shadow
InaccessiblePaths=/var/lib/codex-gateway/cliproxy-auth /var/lib/codex-gateway/cliproxy-auth-archive /etc/cliproxy
LockPersonality=true
