# cloud/systemd/cloudx-shadow-health.service
[Unit]
Description=Cloudx shadow health publisher
After=network-online.target cloudx-shadow-account-state.service
Wants=network-online.target cloudx-shadow-account-state.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/cloudx/cloudx-shadow.env
ExecStart=/opt/cloudx/current/cloudx-cloud.pyz publish-health --json
User=cloudx
Group=cloudx
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/run/cloudx-shadow
